# xware
#!/bin/sh /etc/rc.common

START=90

start() {
	/root/xware/portal
	sleep 3
	curl http://localhost:9000
	exit $?
}

stop() {
	/root/xware/portal -s
}
